# config/autostart/light-locker.desktop
[Desktop Entry]
Type=Application
Name=Screen Locker
Exec=light-locker
TryExec=light-locker
NotShowIn=X-DWM;
