# config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Name=picom
Exec=picom
TryExec=picom
NotShowIn=X-DWM;
