# config/autostart/polkit-mate-authentication-agent-1.desktop
[Desktop Entry]
Type=Application
Name=PolicyKit Authentication Agent
Exec=dwm-polkit
TryExec=dwm-polkit
NotShowIn=X-DWM;
